# systemd/connector-control-api.service
[Unit]
Description=NiceLabel Connector Control API
After=network.target

[Service]
Type=simple
User=nlconnector
WorkingDirectory=/opt/nl-connector/app
ExecStart=/opt/nl-connector/app/.venv/bin/python /opt/nl-connector/app/control_api.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
