# systemd/selector.service
[Unit]
Description=NiceLabel Selector (CalcMenu to Supabase)
After=network.target

[Service]
Type=oneshot
User=nlconnector
WorkingDirectory=/opt/nl-connector/app
ExecStart=/opt/nl-connector/app/.venv/bin/python /opt/nl-connector/app/selector.py
