# systemd/selector.timer
[Unit]
Description=Run NiceLabel Selector Daily at 6AM

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
